# deploy/dashboard-franquicia-web.service
[Unit]
Description=Dashboard Franquicia Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/dashboard-franquicia/frontend
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/npm run start -- -p 3001 -H 127.0.0.1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
